# deploy/systemd/yangyang-image-generator.service.example
[Unit]
Description=YANGYANGAPI Image Workspace
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/yangyang-image-generator/app
EnvironmentFile=/opt/yangyang-image-generator/.env
Environment=DATA_DIR=/opt/yangyang-image-generator/data
ExecStart=/opt/yangyang-image-generator/.venv/bin/gunicorn -b 127.0.0.1:3012 --workers 2 --timeout 660 server:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
